# opkg-feed/installer/openwrt-sdxlemur-skel/usr/lib/opkg/info/busybox-selinux.list
/bin/dmesg
/bin/printenv
/usr/bin/id
/bin/login
/usr/bin/realpath
/bin/rm
/sbin/ip
/usr/sbin/sestatus
/bin/sh
/usr/bin/test
/usr/sbin/load_policy
/bin/zcat
/bin/fgrep
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/usr/sbin/getsebool
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/usr/sbin/setsebool
/usr/sbin/ntpd
/usr/bin/time
/bin/vi
/bin/pwd
/usr/bin/nc
/bin/chown
/etc/init.d/cron
/bin/mount
/bin/sync
/usr/bin/runcon
/usr/bin/dirname
/sbin/hwclock
/usr/sbin/getenforce
/sbin/udhcpc
/bin/ping6
/usr/bin/basename
/usr/bin/chcon
/sbin/restorecon
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/sbin/crond
/usr/bin/head
/bin/gunzip
/bin/netstat
/bin/busybox
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/bin/kill
/bin/mktemp
/usr/bin/tr
/bin/mknod
/bin/rmdir
/usr/bin/unzip
/usr/bin/uniq
/usr/bin/tee
/usr/bin/flock
/lib/upgrade/keep.d/busybox-selinux
/usr/bin/hexdump
/etc/init.d/sysntpd
/bin/traceroute6
/usr/bin/sha256sum
/bin/date
/usr/bin/wc
/usr/bin/[
/usr/bin/nslookup
/usr/bin/pgrep
/bin/umount
/sbin/swapon
/bin/netmsg
/usr/bin/tail
/usr/sbin/ntpd-hotplug
/usr/bin/mkfifo
/bin/chgrp
/sbin/swapoff
/bin/cp
/etc/syslog.conf
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/bin/df
/bin/false
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/bin/traceroute
/usr/bin/cmp
/usr/bin/reset
/usr/bin/fuser
/usr/share/acl.d/ntpd.json
/sbin/start-stop-daemon
/etc/capabilities/ntpd.json
/usr/sbin/setenforce
/usr/sbin/brctl
/bin/cat
/usr/sbin/arping
/bin/passwd
/bin/touch
/sbin/mkswap
/bin/tar
/usr/bin/yes
/usr/bin/top
/sbin/switch_root
/usr/bin/expr
/usr/bin/env
/sbin/sysctl
/bin/sed
/usr/bin/logger
/bin/chmod
/usr/bin/awk
/bin/uname
/usr/bin/which
/usr/bin/find
/usr/sbin/matchpathcon
/bin/ln
/bin/ls
/sbin/pivot_root
/usr/bin/nohup
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/sbin/selinuxenabled
/usr/bin/crontab
/usr/bin/printf
/bin/ping
/sbin/setfiles
/usr/bin/uptime
/bin/fsync
/bin/pidof
/bin/egrep
/usr/bin/readlink
/usr/bin/strings
/sbin/arp
/bin/ps
/usr/bin/xargs
/bin/gzip
