# opkg-feed/installer/openwrt-sdxlemur-skel/data/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
# avoid auto reboot by twice type ctrl+c
#::askconsole:/usr/libexec/login.sh
console::respawn:/usr/libexec/login.sh